# scripts/crontab
# RMS Defragmenter Cron Jobs
# Run defragmentation analysis daily at 2 AM using CLI mode
0 2 * * * appuser cd /app && python start.py >> /app/logs/defrag_cron.log 2>&1

# Empty line required at end of crontab
